# systemd/tweakerPatch.timer
[Unit]
Description=Hammond On Electrix Tweaker Patch timer

[Timer]
OnUnitActiveSec=2s
OnBootSec=2s

[Install]
WantedBy=timers.target
